# deploy/systemd/ollama-bridge-healthcheck.service
[Unit]
Description=Memphis v4 Ollama Bridge Healthcheck + Auto-Recovery
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/home/memphis_ai_brain_on_chain/memphis-v4
ExecStart=/usr/bin/env bash /home/memphis_ai_brain_on_chain/memphis-v4/scripts/ollama-bridge-healthcheck.sh
